# ops/kontur-sync.service
[Unit]
Description=Kontur scheduled content connector sync and freshness check
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/kontur
Environment=COMPOSE_PROJECT_NAME=kontur
ExecStart=/usr/bin/flock -n /run/kontur-sync.lock /usr/bin/docker compose run --rm --no-deps app python -m kontur.cli automation run
TimeoutStartSec=45min
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6
